# vendorbuild/mikrotik/crontab
0 */12 * * * pathvector generate >> /var/log/cron.log 2>&1
